# astradb-core.service
[Unit]
Description=AxiomDB gateway Rust control plane
After=network.target postgresql.service

[Service]
Type=simple
User=opsdc
Group=opsdc
WorkingDirectory=/home/opsdc
EnvironmentFile=/home/opsdc/.creds/astradb-core.env
ExecStart=/usr/local/bin/axiomdb-gateway
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/opsdc

[Install]
WantedBy=multi-user.target
